# docker/files/services/openwebrx/run
#!/usr/bin/execlineb -P
/opt/openwebrx/docker/scripts/run.sh