# docker/files/services/codecserver/run
#!/command/execlineb -P
/usr/local/bin/codecserver