# docker/files/services/codecserver/run
#!/usr/bin/execlineb -P
/usr/local/bin/codecserver