# systemd/openwebrx.service
[Unit]
Description=OpenWebRX WebSDR receiver

[Service]
Type=simple
User=openwebrx
Group=openwebrx
ExecStart=/usr/bin/openwebrx
Restart=always

[Install]
WantedBy=multi-user.target
